# crontab
*/30 * * * * python3 /app/main.py >> /app/logfile